# etc/debian/rib-ribs/files/ribsd.service
[Unit]
Description=Beluganos RIB sync service
Wants=nlad.service beluganos.service
After=syslog.target network.target nlad.service
ConditionPathExists=/etc/beluganos/ribxd.conf

[Service]
Type=simple
<<<<<<< HEAD
<<<<<<< HEAD
<<<<<<< HEAD
ExecStart=/usr/bin/ribsd -config /etc/beluganos/ribxd.conf
=======
ExecStart=/usr/bin/ribsd --config-file /etc/beluganos/ribxd.conf
>>>>>>> develop
=======
ExecStart=/usr/bin/ribsd --config-file /etc/beluganos/ribxd.conf
>>>>>>> develop
=======
ExecStart=/usr/bin/ribsd --config-file /etc/beluganos/ribxd.conf
>>>>>>> develop
User=beluganos
Group=beluganos
Restart=on-abort

[Install]
WantedBy=network.target
