# etc/ztp/debian/postinst
#! /bin/bash
# -*- coding: utf-8; mode: shell-script -*-

# systemctl daemon-reload
# systemctl enable ztp-init
